# dehydrated.cron
#!/bin/bash 

cd /webapp
bash -c "echo $(date) >> hello"
./dehydrated -c -f letsencrypt.conf
